# scripts/a11y-remediate.service
[Unit]
Description=A11y Remediation App
After=network.target

[Service]
Type=exec
User=ubuntu
WorkingDirectory=/home/ubuntu/a11y-remediate
EnvironmentFile=/home/ubuntu/a11y-remediate/.env.production
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk-arm64
Environment=SECURE_COOKIES=true
ExecStart=/home/ubuntu/a11y-remediate/.venv/bin/uvicorn src.web.app:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
